# superctrl.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.superctrl.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/superctrl</string>
        <string>daemon</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/superctrl.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/superctrl.error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>OPENAI_API_KEY</key>
        <string>REPLACE_WITH_YOUR_API_KEY</string>
    </dict>
</dict>
</plist>
